# spring-cloud-kubernetes-demo.ini
[program:user-service-provider]
command=java -jar /usr/local/project/user-service-provider.jar --spring.profiles.active=test

[program:user-service-invoker]
command=java -jar /usr/local/project/user-service-invoker.jar --spring.profiles.active=test